# infra/cron/crontab
# momentum-etf VM 배치 스케줄 (타임존: Asia/Seoul)
#
# ⚠️ 경고: `crontab <파일>` 로 직접 설치하지 말 것.
#    이 명령은 사용자 crontab 전체를 덮어써서 같은 VM 에서 돌아가는 다른 앱
#    (예: leverage-switching) 의 cron 항목을 삭제한다.
#    반드시 아래 install.sh 로만 설치/재설치 한다. install.sh 는 자기
#    마커(# >>> momentum-etf >>> ~ # <<< momentum-etf <<<) 블록만 교체하여
#    다른 앱 항목을 보존하고, 여러 번 실행해도 idempotent 하다.
#
# 설치 / 재설치:
#   bash /home/ubuntu/apps/momentum-etf/infra/cron/install.sh
# 제거(자기 블록만):
#   bash /home/ubuntu/apps/momentum-etf/infra/cron/install.sh --uninstall
# 확인:
#   crontab -l
#
# 동작 원리:
#   호스트 cron 이 docker compose exec 로 fastapi_app 컨테이너 안의
#   infra/cron/run_batch.py 를 호출합니다. 래퍼가 스크립트를 실행하고
#   실패 결과를 Slack 으로 전송합니다. (성공 알림은 비활성화)
#
# 로그:
#   logs/cron/<job>.log 에 append 됩니다.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
APP_DIR=/home/ubuntu/apps/momentum-etf

# 모든 배치는 캘린더 평일(월~금, KST) 에만 실행한다.
# cron 의 요일 필드: 1=월, 2=화, ..., 5=금

# 평일 매시간 정각 - 주가 캐시 갱신
0 * * * 1-5 cd $APP_DIR && mkdir -p logs/cron && docker compose exec -T fastapi_app python infra/cron/run_batch.py cache_refresh python scripts/stock_price_cache_updater.py >> logs/cron/cache_refresh.log 2>&1

# 평일 09:00 KST - 종목 메타데이터 업데이트
0 9 * * 1-5 cd $APP_DIR && mkdir -p logs/cron && docker compose exec -T fastapi_app python infra/cron/run_batch.py metadata_updater python scripts/stock_meta_cache_updater.py >> logs/cron/metadata_updater.log 2>&1

# 평일 09:30, 16:30 KST - 자산 요약
30 9,16 * * 1-5 cd $APP_DIR && mkdir -p logs/cron && docker compose exec -T fastapi_app python infra/cron/run_batch.py asset_summary python scripts/slack_asset_summary.py >> logs/cron/asset_summary.log 2>&1

# 평일 07:00 KST - 장 시간 분석
0 7 * * 1-5 cd $APP_DIR && mkdir -p logs/cron && docker compose exec -T fastapi_app python infra/cron/run_batch.py market_hours_analysis python scripts/analyze_market_hours.py >> logs/cron/market_hours_analysis.log 2>&1
